# systemd/afd.service
[Unit]
Description=Auto delete files


[Service]
ExecStart=afdd %h
CPUWeight=1
IOWeight=1

[Install]
WantedBy=multi-user.target
