# systemd/afd.service
[Unit]
Description=Auto delete files


[Service]
ExecStart=afdd %h

[Install]
WantedBy=multi-user.target
